# learning.service
[Unit]
Description=Learning Navigator Web Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/repos/learning
ExecStart=/home/pi/.local/bin/uv run python server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
