# raspberry_pi_setup/learning.service
[Unit]
Description=Learning Navigator Web Server
After=network.target

[Service]
Type=simple
User=jayab867
WorkingDirectory=/home/jayab867/repos/learning
ExecStart=/home/jayab867/.local/bin/uv run python server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
